# container_files/etc/supervisor/conf.d/gunicorn.conf
[program:gunicorn]
directory=/server
command=gunicorn -b 127.0.0.1:8000 -w 4 %(ENV_APP_NAME)s.wsgi:application --log-file=/opt/logs/gunicorn/gunicorn.log
autostart=true
autorestart=true
